# IoTGoat-raspberry-pi2.img_extracted/rootfs/usr/lib/opkg/info/shellback.list
/etc/init.d/shellback
/usr/bin/shellback
/etc/config/shellback
